# deploy/supervisord.conf
; HealthDesk AI supervisord config (TRD §12.3).
; Runs the FastAPI gateway and the LiveKit voice worker as siblings;
; a crash in either is logged and the container exits 1 so the
; orchestrator restarts it.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

[program:gateway]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
exitcodes=0

[program:voice]
command=python -m app.voice.agent_worker start
autostart=%(ENV_HEALTHDESK_VOICE)s
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
exitcodes=0

[eventlistener:exit_on_any_fatal]
command=bash -c "printf 'READY\n' && while read line; do echo \"$line\" >&2; kill -SIGTERM $PPID; done"
events=PROCESS_STATE_FATAL
